# systemd/corax_orchestrator.service
[Unit]
Description=Corax Autonomous Orchestrator Service
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/Crypto-MCP-Server---by-Corax-CoLAB
EnvironmentFile=%h/Crypto-MCP-Server---by-Corax-CoLAB/.env
ExecStart=/usr/bin/python3 autonomous_orchestrator.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
